# sys-fs/sanoid/files/sanoid.cron
TZ=UTC
*/15 * * * * root flock -n /var/lock/sanoid-cron.lock sanoid --force-update --cron
