# sys-fs/sanoid/files/sanoid.cron
TZ=UTC
*/15 * * * * root flock -n /var/lock/sanoid-cron-take.lock sanoid --take-snapshots
*/15 * * * * root flock -n /var/lock/sanoid-cron-prune.lock sanoid --prune-snapshots
