# app-misc/logiops/files/logid
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="logid"
description="Logitech Configuration Daemon"
command=/usr/bin/logid
command_arg="${logid_args}"
supervisor="supervise-daemon"

depend() {
	need udev
}
